# skel/etc/teamster-backup.timer.tmpl
[Unit]
Description=Teamster data backup (every __SCHEDULE__)

[Timer]
OnBootSec=5min
OnUnitActiveSec=__SCHEDULE__
AccuracySec=1min

[Install]
WantedBy=timers.target
